# files/etc/uci-defaults/90_storage_config
. /lib/functions/touchstar.sh

board=$(boardname_get)

logger -t storage-config "Applying storage configuration [start]"

case "$board" in
ts-vh401)
		# Set fstab defaults
  uci add fstab mount
  uci set fstab.@mount[0].device="/dev/sda1"
  uci set fstab.@mount[0].options="rw,sync"
  uci set fstab.@mount[0].target="/mnt/storage"
  uci set fstab.@mount[0].enabled_fsck="0"
  uci set fstab.@mount[0].enabled="1"
  uci commit fstab
 	;;
esac

logger -t storage-config "Applying storage configuration [end]"
